# crontab
SHELL=/bin/bash

* * * * * root echo "hello at $(date +'\%Y-\%m-\%d \%H:\%M:\%S')" >> /var/log/cron.log 2>&1

